# xmlmesh/server/DEBIAN/postinst
#!/bin/bash
set -e

if [ $1 = "configure" -o $1 = "reconfigure" ]; then
    update-rc.d ot-xmlmesh defaults 19 21
    /etc/init.d/ot-xmlmesh start
fi
