# xmlmesh/server/DEBIAN/prerm
#!/bin/bash
set -e

if [ $1 = "remove" -o $1 = "upgrade" -o $1 = "purge" ]; then
  /etc/init.d/ot-xmlmesh stop
  update-rc.d -f ot-xmlmesh remove > /dev/null 2>&1
fi

